# anpr.service
[Unit]
Description=ANPR Flask Application
After=network.target

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/opt/anpr
Environment=PATH=/opt/anpr/venv/bin
Environment=FLASK_ENV=production
Environment=PYTHONPATH=/opt/anpr
ExecStart=/opt/anpr/venv/bin/gunicorn --bind 127.0.0.1:8000 --workers 4 --timeout 120 --access-logfile - --error-logfile - plate_recognition.app:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
